# etc/udev/rules.d/90-fbcondecor.rules
KERNEL=="fb[0-9]", ACTION!="remove", RUN+="/sbin/splash-term --auto"
