# etc/udev/rules.d/92-console-font.rules
KERNEL=="fb[0-9]", ACTION!="remove", RUN+="/bin/live/set-console-width --auto"
